# scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/slt/supervisord.log
pidfile=/var/run/slt/supervisord.pid
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/slt/nginx-stdout.log
stderr_logfile=/var/log/slt/nginx-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

[program:gunicorn]
command=gunicorn main:app -c /app/gunicorn.conf.py
directory=/app
user=slt
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/slt/gunicorn-stdout.log
stderr_logfile=/var/log/slt/gunicorn-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    SECRET_KEY="%(ENV_SECRET_KEY)s",
    LOG_LEVEL="%(ENV_LOG_LEVEL)s",
    DISCORD_CLIENT_ID="%(ENV_DISCORD_CLIENT_ID)s",
    DISCORD_CLIENT_SECRET="%(ENV_DISCORD_CLIENT_SECRET)s",
    DISCORD_REDIRECT_URI="%(ENV_DISCORD_REDIRECT_URI)s",
    PUBLIC_BASE_URL="%(ENV_PUBLIC_BASE_URL)s",
    POSTMARK_SERVER_TOKEN="%(ENV_POSTMARK_SERVER_TOKEN)s",
    POSTMARK_FROM_EMAIL="%(ENV_POSTMARK_FROM_EMAIL)s",
    POSTMARK_MESSAGE_STREAM="%(ENV_POSTMARK_MESSAGE_STREAM)s"
[group:slt]
programs=nginx,gunicorn
priority=999
